# personal-workstation/mkdeb.src/files/lib/systemd/user/mpd.socket
[Unit]
Description=Music Player Daemon socket

[Socket]
ListenStream=%t/mpd/socket

[Install]
WantedBy=sockets.target
